# ops/systemd/scheduler-replay-report.timer
[Unit]
Description=Nightly Scheduler Replay Calibration Timer

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
Persistent=true
RandomizedDelaySec=10m
Unit=scheduler-replay-report.service

[Install]
WantedBy=timers.target
